# supervisord/conf.d/xvfb.conf
[program:xvfb]
command=Xvfb :0 -screen 0 "%(ENV_DISPLAY_DIMENSIONS)s"x24 -listen tcp -ac
autorestart=true